# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/data/live_logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/tmp/supervisord.pid
user=root

; ── ALL BOTS: Single unified process with shared data layer ─────
; Runs scalp_v2, vol_scalp_v1, production, trident, lotto_yolo,
; lotto_prod, and zeus in one async event loop.
; ONE Polygon WebSocket (delayed stocks) + ONE ORATS cache.
; Accounts: 1-6 (6 separate Alpaca paper accounts)
[program:unified_options]
command=python scripts/run_unified_options.py
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=10
stopwaitsecs=30
stdout_logfile=/app/data/live_logs/unified_options/supervisor.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_redirect_stderr=true
environment=PYTHONPATH="/app/src",PYTHONUNBUFFERED="1",TZ="America/New_York"

; ── Health check HTTP server (for Fly.io) ────────────────────────
[program:healthcheck]
command=python deploy/healthcheck.py
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_redirect_stderr=true
environment=PYTHONPATH="/app/src",PYTHONUNBUFFERED="1"
priority=1

; ── Group ────────────────────────────────────────────────────────
[group:bots]
programs=unified_options
priority=10
